# systemd/wallpaper-brain.timer
[Unit]
Description=Wallpaper brain periodic cleanup and recommendation

[Timer]
OnCalendar=*:0/20
Persistent=true

[Install]
WantedBy=timers.target
